# systemd/ydotoold.service.d/flashpaste-socket.conf
# Drop-in patch for ydotoold.service to fix the socket-path bug in
# ydotool 0.1.8 (Ubuntu 24.04). The daemon ignores --socket-path and
# always listens on /tmp/.ydotool_socket; flashpaste expects it at
# $XDG_RUNTIME_DIR/.ydotool_socket. Symlink bridges them.
[Service]
ExecStartPost=/bin/ln -sf /tmp/.ydotool_socket %t/.ydotool_socket
ExecStopPost=/bin/rm -f %t/.ydotool_socket
